# dev/docker/dionaea/supervisor/rebuild.conf
# This file is part of the dionaea honeypot
#
# SPDX-FileCopyrightText: none
#
# SPDX-License-Identifier: CC0-1.0

[program:rebuild]
command=bash -c 'while inotifywait -q -r -e create,delete,modify,move,attrib --exclude "/\." /code/; do cd /code && /code/dev/docker/dionaea/build.sh && supervisorctl restart dionaea; done'
autostart=true
autorestart=true

stdout_logfile=/var/log/dionaea_dev/rebuild.log
stdout_logfile_maxbytes=0

redirect_stderr=true
